# bash_scripts/services_usr_scripts/montar_drive.service
[Unit]
Description=Montar Google Drive al inicio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/montar_drive.sh
User=gerar_kde
Group=gerar_kde
Restart=on-failure
RestartSec=15
TimeoutSec=60
Environment=PATH=/usr/bin:/bin

# Opcional: limita cuánto puede fallar antes de dejar de intentar
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
